# files/etc/s6/fr24feed/finish
#!/bin/sh

kill "$(ps aux | grep '[f]r24feed' | grep -v 's6' | awk '{print $1}')"
