# files/etc/s6/piaware/finish
#!/bin/sh

kill "$(ps aux | grep '[p]iaware' | grep -v 's6' | awk '{print $1}')"
