# files/etc/s6/pfclient/finish
#!/bin/sh

kill "$(ps aux | grep '[p]fclient' | grep -v 's6' | awk '{print $1}')"
