# files/etc/s6/pfclient/run
#!/bin/sh

exec /usr/bin/pfclient --config_path=/etc/pfclient.json --log_path=/var/log
